# etc/systemd/system/pacman-snapshot.service
[Unit]
Description=Snapshot pacman database

[Service]
Type=oneshot
ExecStart=/bin/bash /etc/systemd/scripts/pacman-snapshot
RemainAfterExit=no
